# examples/config/supervisord/slim-reactor.ini
[program:slim-reactor]
command=/home/slimreactor/examples/strict-psr7-http-server.php 0.0.0.0:40%(process_num)02d
process_name=%(program_name)s-%(process_num)d
numprocs=4
directory=/home/slimreactor/examples
umask=022
user=slimreactor
stdout_logfile=/var/log/supervisord/%(program_name)s-%(process_num)d.log
stderr_logfile=/var/log/supervisord/%(program_name)s-%(process_num)d-error.log
autostart=true
autorestart=true
startretries=3